# deploy/systemd/pdf-server.service
[Unit]
Description=PDF Server API (FastAPI)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/pdf-server
EnvironmentFile=/etc/pdf-server.env
ExecStart=/opt/pdf-server/.venv/bin/uvicorn app.main:app --host ${PDF_SERVER_HOST} --port ${PDF_SERVER_PORT}
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target